# recipes-apps/partition-table-utilities/files/partition-check.service
[Unit]
Description=Check and repair the partition table

[Service]
Type=oneshot
ExecStart=/usr/bin/part-table-check.sh